# files/teslamate-backup.timer
[Unit]
Description=Schedule backups for teslamate postgres

[Timer]
OnBootSec=600
OnUnitActiveSec=3600

[Install]
WantedBy=multi-user.target
